# LichtDomController/Software/supervisord/domControl.conf
[program:domControl]
command=/home/pi/domControl.py
autostart=true
autorestart=true
startretries=10
user=root
directory=/home/pi
